# golang/service/go-service.conf
description "Go test service"
author "Edward Petersen"

start on runlevel [2345]
stop on runlevel [016]

respawn

setuid epeters
setgid epeters
exec /home/epeters/go-workspace/go-service/bin/service
